# client/etc/init/ikn-sync.conf
# ikn-sync

start on desktop-shutdown

task
exec /opt/ikn/sync/start
